# data/rc.local
#!/bin/bash
 /usr/bin/socat tcp-listen:25,reuseaddr,fork tcp:mailhog:1025 &